# root/etc/services.d/fcgiwrap/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/sbin/fcgiwrap -f -s tcp:0.0.0.0:8081
